# root/etc/services.d/darkhttpd/run
#!/usr/bin/with-contenv bash

if [ "$IPV6" = "true" ]; then
	ipv6="--ipv6"
fi

exec \
	s6-setuidgid abc darkhttpd /app/AriaNg --port 80 --no-listing --no-server-id $ipv6 >/dev/null 2>&1
